# P1/supervisord/config.conf
[supervisord]
nodemon=true

[program:apache]
command = /usr/sbin/apache2ctl -D FOREGROUND

[program:node]
command = /usr/bin/node /usr/local/apache2/htdocs/node/ejemplo.js

[supervisorctl]